# squashfs-linux/Linux/lib/systemd/system/dell_logrotate.timer
[Unit]
Description=Log Rotate Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1min

[Install]
WantedBy=multi-user.target
