# squashfs-linux/Linux/lib/systemd/system/dell_logrotate.service
[Unit]
Description=Common Logging Mechanism: Log Rotate

[Service]
Type=oneshot
ExecStart=/bin/dell_logrotate
