# squashfs-linux/Linux/lib/systemd/system/fmgr.socket
[Unit]
Description=feature manager socket

[Socket]
ListenStream=/tmp/fmSocket
Accept=false
Service=fmgr.service

[Install]
WantedBy=sockets.target
