# squashfs-linux/Linux/lib/systemd/system/sdlib_serv.service
[Unit]
Description=socket activation test service
ConditionPathExists=/usr/bin/sdtool_sockserv

[Service]
Type=simple
ControlGroup=memory:/idrac/NC
ExecStart=/bin/sdtool_sockserv

