# squashfs-filesystem/Linux/lib/systemd/system/dbus.socket
[Unit]
After=systemd-tmpfiles-setup.service
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/var/run/dbus/system_bus_socket
